# astralbrowser-indexer.service
[Unit]
Description=generates ***REMOVED*** astralbrowser index

[Service]
Type=oneshot
ExecStart=/opt/***REMOVED***/bin/poc-indexer
User=naggie
