# astralbrowser-indexer.timer
[Unit]
Description=Runs indexer for astralbrowser

[Timer]
OnActiveSec=0
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
