# mounts/mnt-depot.mount
# note these files can be generated by systemd-mount
[Unit]
Description=Central media share
# as we're using DNS hostnames....
After=nss-lookup.target
# Could put dsnet service here as an After= dependency, but this machine does
# not use dsnet for the VPN as it has a direct connection.

[Mount]
What=oxygen.***REMOVED***.io:/srv/media
# must match file name (escaped with systemd-escape!)
Where=/mnt/depot
Options=ro,soft,bg,vers=4
Type=nfs
TimeoutSec=10

[Install]
WantedBy=multi-user.target
